# config/osquery/osquery.conf
{
  "options": {
    "config_plugin": "filesystem",
    "logger_plugin": "filesystem",
    "logger_path": "/var/log/osquery",
    "disable_logging": false,
    "log_result_events": true,
    "logger_min_status": 0,
    "logger_min_stderr": 0,
    "schedule_splay_percent": 10,
    "pidfile": "/var/osquery/osquery.pid",
    "events_expiry": 3600,
    "database_path": "/var/osquery/osquery.db",
    "verbose": false,
    "worker_threads": 2,
    "enable_monitor": true,
    "disable_events": false,
    "disable_audit": false,
    "audit_allow_config": true,
    "audit_allow_sockets": true,
    "host_identifier": "hostname",
    "enable_syslog": false,
    "syslog_pipe_path": "/var/osquery/syslog_pipe",
    "force": false,
    "disable_wd": false,
    "utc": true
  },
  "schedule": {
    "system_info": {
      "query": "SELECT hostname, cpu_brand, physical_memory FROM system_info;",
      "interval": 3600
    },
    "osquery_info": {
      "query": "SELECT version, pid, uuid FROM osquery_info;",
      "interval": 600,
      "snapshot": true
    },
    "processes": {
      "query": "SELECT pid, name, path, cmdline, cwd, root, uid, gid, start_time, parent, pgroup, threads, nice FROM processes;",
      "interval": 600
    },
    "socket_events": {
      "query": "SELECT action, pid, path, auid, success FROM socket_events WHERE success=1;",
      "interval": 60,
      "removed": false
    },
    "process_events": {
      "query": "SELECT auid, cmdline, ctime, cwd, egid, euid, gid, parent, path, pid, time, uid FROM process_events WHERE path NOT IN ('/bin/sed', '/usr/bin/tr', '/bin/gawk', '/bin/date', '/bin/mktemp', '/usr/bin/dirname', '/usr/bin/head', '/usr/bin/jq', '/bin/cut', '/bin/uname', '/bin/basename') and cmdline NOT LIKE '%_key%' AND cmdline NOT LIKE '%secret%';",
      "interval": 60,
      "removed": false
    },
    "file_events": {
      "query": "SELECT target_path, category, time, action FROM file_events;",
      "interval": 60,
      "removed": false
    },
    "network_interfaces": {
      "query": "SELECT interface, address, mask, broadcast, point_to_point FROM interface_addresses;",
      "interval": 600
    },
    "listening_ports": {
      "query": "SELECT pid, port, protocol, family, address FROM listening_ports;",
      "interval": 600
    },
    "crontab": {
      "query": "SELECT * FROM crontab;",
      "interval": 3600
    },
    "iptables": {
      "query": "SELECT filter_name, chain, policy, target, protocol, src_port, dst_port, src_ip, dst_ip FROM iptables;",
      "interval": 600
    },
    "kernel_modules": {
      "query": "SELECT name, size, used_by, status, address FROM kernel_modules;",
      "interval": 600
    },
    "last": {
      "query": "SELECT username, tty, pid, type, time, host FROM last WHERE time > ((SELECT unix_time FROM time) - 3600);",
      "interval": 600
    },
    "mounts": {
      "query": "SELECT device, device_alias, path, type, blocks_size FROM mounts;",
      "interval": 3600
    },
    "shell_history": {
      "query": "SELECT uid, time, command, history_file FROM shell_history WHERE time > ((SELECT unix_time FROM time) - 86400);",
      "interval": 3600
    },
    "startup_items": {
      "query": "SELECT name, path, args, type, source, status, username FROM startup_items;",
      "interval": 3600
    },
    "suid_bin": {
      "query": "SELECT * FROM suid_bin;",
      "interval": 3600
    },
    "system_controls": {
      "query": "SELECT * FROM system_controls WHERE oid LIKE '%.net.%' OR oid LIKE '%.kernel.%';",
      "interval": 3600
    },
    "installed_applications": {
      "query": "SELECT name, path, bundle_version FROM apps;",
      "interval": 3600
    },
    "kernel_info": {
      "query": "SELECT version, arguments FROM kernel_info;",
      "interval": 3600
    },
    "etc_hosts": {
      "query": "SELECT * FROM etc_hosts;",
      "interval": 3600
    },
    "launchd": {
      "query": "SELECT * FROM launchd;",
      "interval": 3600
    },
    "users": {
      "query": "SELECT uid, gid, uid_signed, gid_signed, username, description, directory, shell, uuid FROM users;",
      "interval": 3600
    },
    "groups": {
      "query": "SELECT * FROM groups;",
      "interval": 3600
    },
    "docker_containers": {
      "query": "SELECT id, name, image, image_id, command, created, state, status FROM docker_containers;",
      "interval": 600
    },
    "docker_images": {
      "query": "SELECT id, created, size_bytes, tags FROM docker_images;",
      "interval": 3600
    },
    "docker_networks": {
      "query": "SELECT id, name, driver, created, options FROM docker_networks;",
      "interval": 3600
    },
    "yara_events": {
      "query": "SELECT target_path, category, action, time, strings, tags, count FROM yara_events;",
      "interval": 60,
      "removed": false
    },
    "hardware_events": {
      "query": "SELECT * FROM hardware_events;",
      "interval": 60,
      "removed": false
    },
    "file_integrity_monitoring": {
      "query": "SELECT target_path, md5, sha1, sha256, time FROM hash WHERE path IN ('/etc/passwd', '/etc/shadow', '/etc/sudoers', '/etc/hosts', '/etc/hostname', '/home/ubuntu/cryb-platform/.env');",
      "interval": 300
    },
    "suspicious_processes": {
      "query": "SELECT pid, name, path, cmdline, cwd, uid, gid, start_time, parent FROM processes WHERE name IN ('nc', 'netcat', 'wget', 'curl') OR cmdline LIKE '%bash -i%' OR cmdline LIKE '%sh -i%' OR cmdline LIKE '%/dev/tcp%' OR cmdline LIKE '%base64 -d%';",
      "interval": 60
    },
    "privilege_escalation": {
      "query": "SELECT pid, name, path, cmdline, uid, euid, gid, egid FROM processes WHERE (uid != euid OR gid != egid) AND euid = 0;",
      "interval": 120
    },
    "network_connections": {
      "query": "SELECT pid, fd, socket, family, protocol, local_address, remote_address, local_port, remote_port, state FROM process_open_sockets WHERE remote_address != '' AND remote_address != '127.0.0.1' AND remote_address NOT LIKE '192.168.%' AND remote_address NOT LIKE '10.%' AND remote_address NOT LIKE '172.%';",
      "interval": 300
    }
  },
  "decorators": {
    "load": [
      "SELECT uuid AS host_uuid FROM system_info;",
      "SELECT user AS username FROM logged_in_users ORDER BY time DESC LIMIT 1;"
    ]
  },
  "packs": {
    "incident-response": "/usr/share/osquery/packs/incident-response.conf",
    "vulnerability-assessment": "/usr/share/osquery/packs/vuln-assessment.conf",
    "ossec-rootkit": "/usr/share/osquery/packs/ossec-rootkit.conf"
  },
  "file_paths": {
    "configuration": [
      "/etc/passwd",
      "/etc/shadow",
      "/etc/sudoers",
      "/etc/hosts",
      "/etc/hostname",
      "/etc/crontab",
      "/etc/ssh/sshd_config",
      "/home/ubuntu/cryb-platform/.env",
      "/home/ubuntu/cryb-platform/docker-compose*.yml",
      "/etc/nginx/sites-enabled/",
      "/etc/ssl/certs/"
    ],
    "binaries": [
      "/usr/bin",
      "/usr/sbin",
      "/bin",
      "/sbin",
      "/usr/local/bin",
      "/usr/local/sbin"
    ],
    "logs": [
      "/var/log/auth.log",
      "/var/log/syslog",
      "/var/log/kern.log",
      "/var/log/cron.log",
      "/var/log/nginx/",
      "/home/ubuntu/cryb-platform/logs/"
    ]
  },
  "events": {
    "disable_subscribers": [
      "user_events"
    ]
  },
  "yara": {
    "signatures": {
      "sig_group_1": [
        "/tmp/malware.yar",
        "/tmp/webshell.yar"
      ]
    },
    "file_paths": {
      "sig_group_1": [
        "/home/ubuntu/cryb-platform/uploads/",
        "/tmp/",
        "/var/tmp/",
        "/dev/shm/"
      ]
    }
  }
}